# deploy/systemd/invoice-frontend.service.in
[Unit]
Description=Invoice System Frontend (Vite Preview)
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__PROJECT_ROOT__/frontend
ExecStart=__PROJECT_ROOT__/deploy/start_frontend.sh
Restart=always
RestartSec=3
Environment=NODE_ENV=production
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
