# deploy/systemd/invoice-backend.service.in
[Unit]
Description=Invoice System Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__PROJECT_ROOT__/backend
ExecStart=__PROJECT_ROOT__/deploy/start_backend.sh
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
